# debian/postinst
#!/bin/sh
# postinst script for ndlrdeposit
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        webserver=apache2
        if [ -d /etc/$webserver/conf.d ] && [ ! -e /etc/$webserver/conf.d/ndlrdeposit.conf ]; then
            ln -s ../../ndlrdeposit/apache.conf /etc/$webserver/conf.d/ndlrdeposit.conf
        fi
        if [ -f /etc/init.d/$webserver ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d $webserver reload 3>/dev/null || true
            else
                /etc/init.d/$webserver reload 3>/dev/null || true
            fi
        fi

	datadir=/var/lib/ndlrdeposit
        if [ -d $datadir ]; then
            chown -R www-data:www-data $datadir
            chmod 0750 $datadir
        fi
    ;;


    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
